# systemd/kbd-backlight-auto.service
[Unit]
Description=Auto keyboard backlight based on ambient light sensor
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/lars/dotfiles-arch/scripts/kbd-backlight-auto.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
